# lesson5/auto/watcher.service
[Unit]
Description=My watcher service

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/watcher
ExecStart=/opt/watcher.sh $WORD $LOG